# systemd/coding-agent-poll@.service
[Unit]
Description=coding-agent-workflow poll (instance: %i)
Documentation=https://github.com/luosky/coding-agent-work-loop
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# instance name (%i) 是 setup.sh 选定的 project key
# EnvironmentFile 在 ~/.config/coding-agent-workflow/<key>.conf 注入：
#   PROJECT_ROOT=...
#   CODING_AGENT_CONFIG=...
#   PATH=...
EnvironmentFile=%h/.config/coding-agent-workflow/%i.conf
# Skill 规范位置：~/.agents/skills/<name>/。各 agent runtime（如 Claude Code 的 ~/.claude/skills/）
# 都建议 symlink 过来，单源真相在这里。
WorkingDirectory=%h/.agents/skills/coding-agent-workflow
ExecStart=/usr/bin/env bash %h/.agents/skills/coding-agent-workflow/scripts/agent-poll.sh

[Install]
WantedBy=default.target
